# Pi433Daemon.conf
description "Pi433Daemon"
author "Matthias Fischer"
start on runlevel [2345]    

pre-start script
  echo "[`date`] Pi433Daemon Starting" >> /var/log/Pi433Daemon.log
end script

exec /bin/sh ___PATH___/run.sh > /dev/null &