# Pi433.service
***Pi433.service***
[Unit]
Description=Pi433-Service
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/flask run --host=0.0.0.0
User=pi
WorkingDirectory=___PATH___
Restart=always
Environment="FLASK_APP=___PATH___/scripts/pi-433.py"

[Install]
WantedBy=multi-user.target